# Route53DynUpdate
#!/sbin/openrc-run

# Catch script Output when in Debug mode
CatchOutput=False

depend() {
    need net
}

checkconfig() {
    # Verify that the configuration file exists
    if [ ! -f /etc/network/Route53DynUpdate.ini ]; then
        eerror "Route53DynUpdate Configuration file /etc/network/Route53DynUpdate.ini doesn't exist"
        return 1
    else
        if grep "Debug=True" /etc/network/Route53DynUpdate.ini &> /dev/null
        then
            CatchOutput=True
        fi
    fi
}

start() {
    checkconfig || return 1
    ebegin "Starting Route53DynUpdate"
    if [ $CatchOutput = True ]
    then
        /opt/bin/Route53DynUpdate.py > /var/log/Route53DynUpdate.stdout.log 2> /var/log/Route53DynUpdate.stderr.log
    else
        /opt/bin/Route53DynUpdate.py
    fi
    eend $?
}

stop() {
    PID=$(cat "/run/Route53DynUpdate.pid" 2>/dev/null)
    if [ -z "${PID}" ]; then
        einfo "Route53DynUpdate not running (no pid file)"
        return 0
    fi

    ebegin "Stopping Route53DynUpdate"
    if ! kill $PID > /dev/null 2>&1; then
        einfo "Could not send SIGTERM to process $PID" >&2
    fi
    eend $?
}

restart() {
    stop
    start
}
